# Weblive/supervisor/crm-queue.conf
; /etc/supervisor/conf.d/crm-queue.conf
;
; 2 Laravel queue worker processes. Restart on crash. Drain on stop.
;
; To restart workers after a code deploy:
;   sudo supervisorctl restart crm-queue:*

[program:crm-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/crm-erp/backend/artisan queue:work redis --queue=default,reminders --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crm-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
